# systemd/moab-hat-flash.service
[Unit]
Description=Flash Project Moab HAT MCU firmware (oneshot, idempotent)
DefaultDependencies=yes
After=local-fs.target systemd-modules-load.service
Before=moab-balance.service
ConditionPathExists=/opt/moab-balance/firmware/v3.bin

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/moab-balance/firmware/flash-hat.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitIntervalSec=300

[Install]
WantedBy=multi-user.target
