# systemd/moab-balance.service
[Unit]
Description=Moab minimal ball-balancer
Documentation=https://github.com/microsoft/moabian
# Wait for the camera and SPI device to be ready.
After=multi-user.target dev-spidev0.0.device dev-video0.device
Wants=dev-spidev0.0.device dev-video0.device
# Make sure the stock Moab service can never race us for the HAT.
Conflicts=moab.service

[Service]
Type=simple
User=pi
Group=pi
SupplementaryGroups=spi gpio video
WorkingDirectory=/opt/moab-balance
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/moab-balance/venv/bin/python -m moab_balance
Restart=on-failure
RestartSec=2
# If the loop exits cleanly (e.g. a controlled shutdown), don't loop forever.
SuccessExitStatus=0

[Install]
WantedBy=multi-user.target
